# paddleocr-daemon.service.example
[Unit]
Description=PaddleOCR Daemon
After=network.target

[Service]
Type=simple
User=paddleocr
Group=paddleocr
EnvironmentFile=-/etc/paddleocr/paddleocr.conf
WorkingDirectory=/opt/paddleocr-daemon
ExecStart=/usr/bin/python3 /opt/paddleocr-daemon/daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
